# project/ry/dockers/php7/conf/supervisord.conf
[supervisord]
nodaemon=true
[program:php]
command=/bin/bash -c "/usr/local/bin/php -c /etc/php.ini"